# packages/mc-server/lib/systemd/system/minecraft@.service
[Unit]
Description=Minecraft Server (%i)
After=network.target
Documentation=https://apt.bulmer.dev

[Service]
Type=simple
User=minecraft
Group=minecraft
WorkingDirectory=/opt/minecraft/%i

ExecStart=/usr/lib/mc/start.sh %i

SuccessExitStatus=0 1
Restart=on-failure
RestartSec=10s
RestartPreventExitStatus=0

TimeoutStartSec=60s
TimeoutStopSec=90s

[Install]
WantedBy=multi-user.target
